# crontab
*/15 * * * *            /inject-speed-metrics.py >> /var/log/speed-metrics.log
